# deb-package/DEBIAN/postinst
#if [ -f /etc/init.d/tikca ]; then
#	rm /etc/init.d/tikca
#fi

#if [ -f /etc/init/tikca.conf ]; then
#    ln -s /etc/init/tikca.conf /etc/init.d/tikca
#
#fi

#if [ ! -d "/var/log/tikca" ]; then
#  install -d -o matterhorn -g nogroup -m755 /var/log/tikca
#  touch /var/log/tikca/tikca.log && chmod 644 /var/log/tikca/tikca.log
#  touch /var/log/tikca/tikca_ingest.log && chmod 644 /var/log/tikca/tikca_ingest.log
#  touch /var/log/tikca/tikca_comm.log && chmod 644 /var/log/tikca/tikca_comm.log
#fi

python3 -m compileall /usr/share/tikca/tikca.py /usr/share/tikca/tikca > /dev/null

update-rc.d tikca defaults
update-rc.d tikca-ingester defaults
service tikca restart
service tikca-ingester restart
