# deb-package/DEBIAN/postinst
#if [ -f /etc/init.d/tikca ]; then
#	rm /etc/init.d/tikca
#fi

#if [ -f /etc/init/tikca.conf ]; then
#    ln -s /etc/init/tikca.conf /etc/init.d/tikca
#
#fi

if [ ! -d "/var/log/tikca" ]; then
  install -d -o root -g root -m777 /var/log/tikca
  touch /var/log/tikca/tikca.log && chmod 666 /var/log/tikca/tikca.log
fi

python3 -m compileall /usr/share/tikca/tikca.py /usr/share/tikca/tikca > /dev/null

update-rc.d tikca defaults
service tikca start